# roles/cfg/files/user/desktop/.config/autostart/workrave.desktop
[Desktop Entry]
Exec=workrave
Comment=RSI prevention tool
Name=workrave
Type=Application
